# services/mhddos/run
#!/bin/sh
exec 2>&1
exec "$SCRIPT_DIR"/bin/mhddos_proxy_linux 2>&1
