# services/mhddos/run
#!/bin/sh
exec 2>&1
exec /opt/itarmy/bin/mhddos_proxy_linux