# services/distress.service
[Unit]
Description=distress
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=on-failure
RestartSec=10
ExecStart=/opt/itarmy/bin/distress
WorkingDirectory=/opt/itarmy/bin
SyslogIdentifier=distress
SyslogFacility=daemon
StandardOutput=append:/var/log/adss.log
StandardError=append:/var/log/adss.log

[Install]
WantedBy=multi-user.target