# services/distress/run
#!/bin/sh
exec 2>&1
exec "$SCRIPT_DIR"/bin/distress
