# services/openrc/db1000n
#!/sbin/openrc-run

name="Start db1000n ddos tool"
pidfile="/run/db1000n.pid"
output_log="/var/log/adss.log"
error_log="/var/log/adss.log"
command="/opt/itarmy/bin/db1000n"
command_background="true"
command_args=""