# services/openrc/distress
#!/sbin/openrc-run

name="Start distress ddos tool"
pidfile="/run/distress.pid"
output_log="/var/log/adss.log"
error_log="/var/log/adss.log"
command="/opt/itarmy/bin/distress"
command_background="true"
command_args=""