# services/openrc/mhddos
#!/sbin/openrc-run

name="Start mhddos_proxy_linux ddos tool"
pidfile="/run/mhddos_proxy_linux.pid"
output_log="/var/log/adss.log"
error_log="/var/log/adss.log"
command="/opt/itarmy/bin/mhddos_proxy_linux"
command_background="true"
command_args=""

